# deploy/hack-detector.service
[Unit]
Description=Hack Detector Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=h970513_gmail_com
WorkingDirectory=/opt/hack-detector
EnvironmentFile=/opt/hack-detector/.env
ExecStart=/opt/hack-detector/.venv/bin/python -m src.main listen
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# e2-micro 리소스 보호
MemoryMax=512M

[Install]
WantedBy=multi-user.target
